# profile.d/color.sh
# Enable colors
export CLICOLOR=1

# Normal colors
export AF_COLOR_BLACK='\e[0;30m'
export AF_COLOR_RED='\e[0;31m'
export AF_COLOR_GREEN='\e[0;32m'
export AF_COLOR_YELLOW='\e[0;33m'
export AF_COLOR_BLUE='\e[0;34m'
export AF_COLOR_MAGENTA='\e[0;35m'
export AF_COLOR_CYAN='\e[0;36m'
export AF_COLOR_WHITE='\e[0;37m'

# Bold colors
export AF_COLOR_BOLD_BLACK='\e[1;30m'
export AF_COLOR_BOLD_RED='\e[1;31m'
export AF_COLOR_BOLD_GREEN='\e[1;32m'
export AF_COLOR_BOLD_YELLOW='\e[1;33m'
export AF_COLOR_BOLD_BLUE='\e[1;34m'
export AF_COLOR_BOLD_MAGENTA='\e[1;35m'
export AF_COLOR_BOLD_CYAN='\e[1;36m'
export AF_COLOR_BOLD_WHITE='\e[1;37m'

# Background colors
export AB_COLOR_BLACK='\e[40m'
export AB_COLOR_RED='\e[41m'
export AB_COLOR_GREEN='\e[42m'
export AB_COLOR_YELLOW='\e[43m'
export AB_COLOR_BLUE='\e[44m'
export AB_COLOR_MAGENTA='\e[45m'
export AB_COLOR_CYAN='\e[46m'
export AB_COLOR_WHITE='\e[47m'

export COLOR_RESET='\e[0m'
export COLOR_ALERT=${AF_COLOR_BOLD_WHITE}${AB_COLOR_RED}
